# dot_config/autostart/private_OneDriveGUI.AppImage.desktop
[Desktop Entry]
Exec=/home/yp/Applications/OneDriveGUI.AppImage
Icon=
Name=OneDriveGUI.AppImage
Path=
Terminal=False
Type=Application
